# root/etc/services.d/flexget/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc flexget -c /flexget/config.yml --loglevel info daemon start
